# sparse/usr/lib/systemd/system/nfcd.service.d/10-wait-for-nfc.conf
[Service]
ExecStartPre=/usr/bin/droid/waitforservice init.svc.nfc_hal_service running
